# com.mk.activitylogger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.mk.activitylogger</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/mk/scripts/activitylogger/start_logger.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/mk/scripts/activitylogger</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/mk/scripts/activitylogger/logs/launchd-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/mk/scripts/activitylogger/logs/launchd-stderr.log</string>
</dict>
</plist>
